# turku_agent_rsyncd.rc.d-freebsd
#!/bin/sh

# PROVIDE: turku_agent_rsyncd
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

# Sample /etc/turku-agent/config.d/freebsd.json :
#
# {
#   "lock_dir": "/var/run",
#   "rsyncd_user": "root",
#   "rsyncd_group": "wheel",
#   "rsyncd_service_name": "turku_agent_rsyncd"
# }

. /etc/rc.subr

export PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin

name="turku_agent_rsyncd"
rcvar="turku_agent_rsyncd_enable"
load_rc_config $name

turku_agent_rsyncd_enable=${turku_agent_rsyncd_enable:-"NO"}

pidfile="/var/run/$name.pid"
command="/usr/sbin/daemon"
command_args="-c -f -P $pidfile /usr/bin/env turku-agent-rsyncd-wrapper"
run_rc_command "$1"
