# turku-agent-ping.timer
[Unit]
Description=turku-agent-ping

[Timer]
OnUnitActiveSec=5m
RandomizedDelaySec=5m
OnStartupSec=15m

[Install]
WantedBy=timers.target
