# turku-agent-rsyncd.conf
description "turku rsync daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/rsync --no-detach --daemon --config=/etc/turku-agent/rsyncd.conf
