# turku-update-config.timer
[Unit]
Description=turku-update-config

[Timer]
OnUnitActiveSec=12h
RandomizedDelaySec=12h
OnStartupSec=1h

[Install]
WantedBy=timers.target
