# turku-update-config.service
[Unit]
Description=turku-update-config

[Service]
Type=oneshot
ExecStart=/usr/bin/env turku-update-config
