# turku-agent-ping.service
[Unit]
Description=turku-agent-ping

[Service]
Type=oneshot
ExecStart=/usr/bin/env turku-agent-ping
